# src/zuya-flask/deployment/roles/tokyogo/templates/supervisor.conf
[program:tokyogo]
command=python app.py
directory=/opt/tokyogo
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
