# src/zuya-flask/deployment/roles/flask_app/templates/supervisor.conf
[program:my_app]
command=python app.py
directory=/opt/my_app
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
